# roles/tinc/templates/tinc.service.j2
[Unit]
Description=tinc vpn
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/tincd -n {{ vpn_interface }}


[Install]
WantedBy=multi-user.target
